# dist/master/cron.d/master.cron
*/1 * * * * root /bin/bash /opt/soft/acl-master/sh/master_keepalive.sh >/dev/null 2>&1
