# dist/master/system/acl-master.service
[Unit]
Description=LSB: acl - master framework to control services
After=network.service

[Service]
Type=forking
Restart=always
RestartSec=3
TimeoutSec=1min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=no
ExecStart=/etc/init.d/acl-master start
ExecStop=/etc/init.d/acl-master stop
ExecReload=/etc/init.d/acl-master reload

[Install]
WantedBy=multi-user.target
